# utils/backup.cron
0 */6 * * * bash /opt/rt3/utils/backup_rt3.sh >> /app/backup/cron.log 2>&1
